# includes.container/etc/systemd/system/batterylimit.service
[Unit]
Description=Duo Battery Limit Service

[Service]
ExecStart=/usr/bin/duo bat-limit

[Install]
WantedBy=multi-user.target