# includes.container/etc/systemd/system/screenmanager.service
[Unit]
Description=Zenbook Duo Monitor service

[Service]
ExecStart=/bin/bash /usr/bin/duomonitor

[Install]
WantedBy=multi-user.target